# pkgs/backend/debian/postinst
#!/bin/sh
set -e
systemctl enable bff.service
systemctl start bff.service
