# pkgs/backend/debian/postrm
#!/bin/sh
set -e
systemctl disable bff.service || true
systemctl stop bff.service || true
